# etc/Morrigan/development/parenchym.supervisor.ini
[program:parenchym]
command=/home/dm/myprojects/Parenchym-py33-venv/bin/gunicorn -c /home/dm/myprojects/Parenchym/Parenchym/etc/Morrigan/development/gunicorn_conf.py --paste /home/dm/myprojects/Parenchym/Parenchym/development.ini
directory=/home/dm/myprojects/Parenchym/Parenchym
environment=PYTHONPATH="/home/dm/myprojects/Parenchym-py33-venv/bin/:/home/dm/myprojects/Parenchym-py33-venv/lib/python3.3/site-packages"
user=dm
autostart=true
autorestart=true
priority=991
stopsignal=KILL
